# machines/led/system/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin potato - $TERM

; [Service]
; ExecStart=
; ExecStart=-/sbin/agetty -n -o potato %I
