# machines/krastiva/system/etc/systemd/system/getty@tty1.service.d/skip-username.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- potato' --noclear --skip-login - $TERM
; ExecStart=-/sbin/agetty -n -o potato %I
