# machines/krastiva/system/etc/udev/rules.d/90-libvirt-usb.rules
SUBSYSTEM=="usb",RUN+="/etc/libvirt/hooks/udev/usb"
